# docker/supervisor/site.conf
[program:site]
command=/var/www/html/dmojsite/bin/uwsgi --ini uwsgi.ini
directory=/var/www/html/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
